# rootfs/etc/services.d/mqttwarn/run
#!/usr/bin/with-contenv sh

/usr/bin/python /opt/mqttwarn/mqttwarn.py
